# deploy/systemd/stranger.service
[Unit]
Description=Stranger video chat (Hono)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/stranger
EnvironmentFile=/opt/stranger/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=20
# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
